# systemd-services/ngulia-check-network.service
# system Daemon file
[Unit]
Description=Run 'check-network.sh'

[Service]
Type=simple
ExecStart=/bin/bash /home/pi/check-network.sh
User=pi
Group=pi
