# systemd-services/ngulia-check-network.service
# system Daemon file
[Unit]
Description=Check if network is available

[Service]
Type=simple
ExecStart=/bin/bash /home/pi/check-network.sh

