# systemd-services/ngulia-check-network.timer
[Unit]
Description=test

[Timer]
OnBootSec=2m
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target