# systemd-services/ngulia-check-network.service
# system Daemon file
[Unit]
Description=Run 'check-network.sh'

[Service]
Type=simple
ExecStart=/usr/bin/sudo /bin/bash /home/pi/check-network.sh

